# bootstrap/97-kloud.tpl.rules
KERNEL=="eth*", ACTION=="add", RUN+="%s"
KERNEL=="wlan*", ACTION=="add", RUN+="%s"
